# automation/com.ebuyingstore.verns-estatesales-monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.ebuyingstore.verns-estatesales-monitor</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/mikesmac/Documents/Codex/2026-05-31/project-title-vern-s-website-create/Vern’s website/scripts/run-estate-sales-monitor.sh</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/mikesmac/Documents/Codex/2026-05-31/project-title-vern-s-website-create/Vern’s website</string>
  <key>StartCalendarInterval</key>
  <array>
    <dict>
      <key>Hour</key>
      <integer>9</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>
    <dict>
      <key>Hour</key>
      <integer>17</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>
  </array>
  <key>StandardOutPath</key>
  <string>/Users/mikesmac/Documents/Codex/2026-05-31/project-title-vern-s-website-create/Vern’s website/output/estate-sales-monitor/launchd.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/mikesmac/Documents/Codex/2026-05-31/project-title-vern-s-website-create/Vern’s website/output/estate-sales-monitor/launchd.err.log</string>
</dict>
</plist>
